# charge-mode
#!/sbin/openrc-run

description="This init script starts charge mode script for leste"
name="Mode Control Entity."

start() {
	/usr/bin/charge-mode.sh &
}

stop() {
    echo no-op
}

