# debian/charge-mode.postinst
#!/bin/sh

install -d /etc/runlevels/charge-mode/
rc-update add charge-mode charge-mode
